# backend/deploy/systemd/jinhengtai-backend.service
[Unit]
Description=Jinhengtai Mall FastAPI Service
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/bash /opt/jinhengtai-backend/deploy/scripts/setup.sh
ExecStop=/usr/bin/docker compose -f /opt/jinhengtai-backend/docker-compose.yml down || /usr/bin/docker-compose -f /opt/jinhengtai-backend/docker-compose.yml down
WorkingDirectory=/opt/jinhengtai-backend
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
